# scripts/robot.service
[Unit]
Description=Robot Tank Control Web Application
Documentation=https://github.com/yourusername/robot
After=network.target motion.service
Wants=motion.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/robot
Environment="RACK_ENV=production"
Environment="PATH=/home/pi/.rbenv/shims:/home/pi/.rbenv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="RBENV_ROOT=/home/pi/.rbenv"

# Start the application using rbenv
ExecStart=/home/pi/.rbenv/shims/bundle exec puma -C config/puma.rb

# Restart policy
Restart=always
RestartSec=10

# Resource limits (important for Pi Zero 2W)
MemoryLimit=200M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=robot-tank

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
