# scripts/rpicam-stream.service
[Unit]
Description=Rpicam Video Streaming Service
Documentation=https://www.raspberrypi.com/documentation/computers/camera_software.html
After=network.target

[Service]
Type=simple
User=pi
Group=video
WorkingDirectory=/tmp

# Stream camera using rpicam-vid piped to ffmpeg for HTTP MJPEG
# rpicam-vid outputs to stdout, ffmpeg serves HTTP multipart MJPEG on port 8081
ExecStart=/bin/bash -c '/usr/bin/rpicam-vid \
    --inline \
    --nopreview \
    --timeout 0 \
    --width 640 \
    --height 480 \
    --framerate 15 \
    --codec mjpeg \
    --quality 85 \
    -o - | \
    /usr/bin/ffmpeg \
    -probesize 32 \
    -analyzeduration 0 \
    -fflags nobuffer \
    -f mjpeg \
    -i - \
    -c:v copy \
    -f mpjpeg \
    -listen 1 \
    http://0.0.0.0:8081'

# Restart policy
Restart=always
RestartSec=5

# Resource limits (important for Pi Zero 2W)
MemoryLimit=250M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpicam-stream

[Install]
WantedBy=multi-user.target
